# dotfiles/python/profile.d/python.sh
#!/bin/sh

# set location of initialization file for python repl
export PYTHONSTARTUP="${HOME}/.pythonrc"
